# deb/install/prerm
#!/bin/bash

systemctl stop sage-wagman-watchdog.timer || true
systemctl disable sage-wagman-watchdog.timer || true
