# deb/install/prerm
#!/bin/bash

systemctl stop sage-wagman-watchdog.service || true
systemctl disable sage-wagman-watchdog.service || true
