# deb/install/prerm
#!/bin/bash

systemctl stop waggle-wagman-watchdog.timer || true
systemctl disable waggle-wagman-watchdog.timer || true
