# ROOTFS/etc/systemd/system/sage-wagman-watchdog.service
[Unit]
Description=Sage Wagman Watchdog

[Service]
ExecStart=/usr/bin/sage-wagman-watchdog
RuntimeMaxSec=10sec

[Install]
WantedBy=multi-user.target
